# etc/seedbox/config/user.crontab.default
@reboot cd ~; sleep $[ ($RANDOM % 60 + 10) ];./.rtorrentExecute.php
* * * * * cd ~; ./.rtorrentRestart.php
